# examples/yabsmd.service
[Unit]
Description=Btrfs Snapshot Manager

[Service]
Type=exec
ExecStart=/usr/bin/yabsm start

[Install]
WantedBy=multi-user.target
